# cloudflare-updater/cloudflare-ip-updater.timer
[Unit]
Description=Update Cloudflare IPs weekly
Requires=cloudflare-ip-updater.service

[Timer]
# Run every Sunday at 3 AM
OnCalendar=Sun *-*-* 03:00:00
# Also run 5 minutes after boot (useful for testing)
OnBootSec=5min
# If the system was off during scheduled time, run when it comes back up
Persistent=true

[Install]
WantedBy=timers.target
