# cloudflare-updater/cloudflare-ip-updater.service
[Unit]
Description=Update Cloudflare IP ranges in Traefik config
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=__INSTALL_DIR__/update-cloudflare-ips.py

# Don't restart on failure - let the timer handle scheduling
Restart=no

# Log output for debugging
StandardOutput=journal
StandardError=journal